# images/php-fpm/context/rootfs/etc/profile.d/history.sh
shopt -s histappend
export HISTFILESIZE=1000000
export HISTSIZE=1000000
export HISTCONTROL=ignoreboth
